# deploy/systemd/blog-api.service
[Unit]
Description=Personal Blog API (Spring Boot)
After=network.target mysql.service redis.service
Wants=mysql.service redis.service

[Service]
Type=simple
User=blog
Group=blog
WorkingDirectory=/opt/blog/app/backend

EnvironmentFile=-/opt/blog/config/blog-api.env
Environment=JAVA_OPTS=-Xms256m -Xmx512m

ExecStart=/usr/bin/java \
  ${JAVA_OPTS} \
  -jar /opt/blog/app/backend/target/personal-blog-api-1.0.0.jar \
  --spring.config.additional-location=file:/opt/blog/config/application-prod.yml

Restart=on-failure
RestartSec=5
SuccessExitStatus=143

# 上传目录、日志
ReadWritePaths=/opt/blog/data /opt/blog/app/backend

[Install]
WantedBy=multi-user.target
