# deploy/systemd/blog-api.service
[Unit]
Description=Personal Blog API (Spring Boot)
After=network.target mysql.service redis.service
Wants=mysql.service redis.service

[Service]
Type=simple
User=blog
Group=blog
WorkingDirectory=/opt/blog/backend

EnvironmentFile=-/opt/blog/config/blog-api.env
Environment=JAVA_OPTS=-Xms256m -Xmx512m

ExecStart=/usr/bin/java \
  ${JAVA_OPTS} \
  -jar /opt/blog/backend/target/personal-blog-api-1.0.0.jar \
  --spring.profiles.active=prod

Restart=on-failure
RestartSec=5
SuccessExitStatus=143

ReadWritePaths=/opt/blog/data /opt/blog/backend

[Install]
WantedBy=multi-user.target
